# scripts/10-my.rules
ACTION=="add", SUBSYSTEM=="block" , RUN+="/home/pi/isoinstaller/scripts/installimage.sh '$env{DEVNAME}'"
